# daemon/transmission-daemon.service
[Unit]
Description=Transmission BitTorrent Daemon

[Service]
User=transmission
Type=notify
ExecStart=/usr/bin/transmission-daemon -f --config-dir /var/lib/transmission/.config/transmission --log-error

[Install]
WantedBy=multi-user.target
